# data/supervisord.conf
[supervisord]
user=root
nodaemon=true

[program:vcontrold]
command=vcontrold -n

[program:crond]
command=crond -f -L /dev/stdout
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:polling]
command=/data/polling.sh
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:listening]
command=/data/listening.sh
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:init]
command=/data/init.sh
autorestart=false
startsecs=0
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
